# system/wifi/dnsmasq.service.d/lume-ap.conf
[Unit]
After=lume-wifi-ap-ip.service
Requires=lume-wifi-ap-ip.service
